# profile.d/00_bashrc.sh
[[ -f ~/.bashrc ]] && grep -q "COMMA_VERSION" ~/.bashrc || source ~/.bashrc

