# profile.d/20_declare.sh
# Global variable declarations

# Bash sometimes struggles with $'\n' depending on the context
readonly LF=$'\n'
readonly HT=$'\t'

# Comma profile vars
COMMA_PROFILEDIR="$HOME/.local/share/comma_bash_profile"
COMMA_VERSION="1.0.7"
# %COMMA_DATE_GENERATED%

